# sv/mongodb/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

DBPATH=${DBPATH="/data/db"}
mkdir -p $DBPATH

exec 2>&1
exec /mongodb/bin/mongod --dbpath=$DBPATH $ARGS
