# sv/mongodb/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

mkdir -p /data/db

exec 2>&1
exec /mongodb/bin/mongod $ARGS
